# supervisord.conf
# [supervisord]
# nodaemon=true

# [program:php-fpm]
# command=php-fpm -F
# autorestart=true
# stdout_logfile=/var/log/php-fpm.log
# stderr_logfile=/var/log/php-fpm.err

# [program:nginx]
# command=/usr/sbin/nginx -g "daemon off;"
# autorestart=true
# stdout_logfile=/var/log/nginx.log
# stderr_logfile=/var/log/nginx.err

[supervisord] 
nodaemon=true 
user=root 
logfile=/var/log/supervisor/supervisord.log 
pidfile=/var/run/supervisord.pid 

[program:php-fpm] 
command=php-fpm -F -R 
autorestart=true 
autostart=true 
stdout_logfile=/var/log/supervisor/php-fpm.log 
stderr_logfile=/var/log/supervisor/php-fpm.log 
stdout_logfile_maxbytes=0 
stderr_logfile_maxbytes=0 
user=root 
priority=1 

[program:nginx] 
command=/usr/sbin/nginx -g "daemon off;" 
autorestart=true 
autostart=true 
stdout_logfile=/var/log/supervisor/nginx.log 
stderr_logfile=/var/log/supervisor/nginx.log 
stdout_logfile_maxbytes=0 
stderr_logfile_maxbytes=0 
user=root 
priority=2 
depends_on=php-fpm